# worker/sentiment-worker.service
[Unit]
Description=Sentiment Arbitrage Worker — fetch Reddit posts, score with FinBERT, write to Firestore
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=nico
WorkingDirectory=/home/nico/sentiment-arbitrage/worker
ExecStart=/home/nico/sentiment-arbitrage/worker/venv/bin/python main.py
EnvironmentFile=/home/nico/sentiment-arbitrage/worker/.env
TimeoutStartSec=300
# Restart on failure with backoff (systemd oneshot retry via timer, but this catches crashes)
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
